# cron/cron23/crontab
* * * * * root python3 /code/t2.py >> /var/log/cron.log
# Empty line
